# Common/battery-shutdown.cron
SHELL=/bin/sh
PATH=/etc:/bin:/sbin:/usr/bin:/usr/sbin

* * * * * root /usr/local/sbin/battery-shutdown.sh > /var/log/battery-shutdown 2>&1
